# docker_base/etc/services.d/haproxy-b/finish
#!/usr/bin/with-contenv sh

# TODO: maybe add checks for bad HAProxy exits (like due to bad error)
